# gongkao-tunnel.service
[Unit]
Description=公考助手 Cloudflare 命名隧道（固定网址 gk.gongkaopei2026.click）
After=network-online.target gongkao.service
Wants=network-online.target

[Service]
Type=simple
Environment=NO_PROXY=*
ExecStart=/home/pei-gen-li/.local/bin/cloudflared tunnel run gongkao
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
